# systemd/user/alert-battery.service
[Unit]
Description=Desktop alert warning of low remaining battery

[Service]
Type=oneshot
# Change username and modify path to script as needed
ExecStart=/home/phuong/.config/scripts/alert-battery

[Install]
WantedBy=graphical.target
